# letsencrypt/geoip2-nginx-stats/root/etc/services.d/geoip2influx/run
#!/usr/bin/with-contenv bash

  exec \
        python3 /config/geoip2influx/geoip2influx.py
